# rootfs/home/jellyfin/.config/autostart/cec-bridge.desktop
[Desktop Entry]
Type=Application
Name=CEC Bridge
Exec=/home/jellyfin/bin/cec-jellyfin.sh
X-GNOME-Autostart-enabled=true
Terminal=false
