# rootfs/home/jellyfin/.config/autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=Unclutter
Exec=unclutter --timeout 1 --jitter 4
X-GNOME-Autostart-enabled=true
Terminal=false
